# Linux/Week_5_Homework_Submission/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.RY09H2/crontab installed on Mon May 18 22:43:09 2020)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
#Ansible: Backup
*/2 * * * * /usr/local/bin/trigger_backup
#Ansible: check for rootkits
@daily bash /opt/chkrootkit
@daily journalctl -p 0..2 > /home/sysadmin/Priority_High.txt
